# debian/cron.d
10 4 * * 1	wdb	/usr/bin/wdb-create-partition --next-period
